# dynv6.timer
[Unit]
Description=Dynv6 update timer

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
